# automation/com.trader.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent: local account health check on weekdays after the US close.
  Time is LOCAL (Europe/London). 21:15 London ~= 16:15 ET (just after close).
  Absolute paths required (launchd does not expand ~ / $HOME).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trader.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/snigdhatiwari/trader/automation/run_monitor.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/snigdhatiwari/trader/logs/launchd-monitor.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/snigdhatiwari/trader/logs/launchd-monitor.err</string>
</dict>
</plist>
