# automation/com.trader.trade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent: one daily rebalance pass on weekdays.
  Time is LOCAL (Europe/London). 15:00 London ~= 10:00 ET (mid US session).
  launchd does NOT expand ~ or $HOME, so paths are absolute. If your username is
  not "snigdhatiwari", search/replace the path below.
  NOTE: launchd only fires while the Mac is awake; a missed weekday slot runs once
  on next wake. For 24/7 reliability use a VPS instead (see automation/README.md).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trader.trade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/snigdhatiwari/trader/automation/run_trade.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/snigdhatiwari/trader/logs/launchd-trade.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/snigdhatiwari/trader/logs/launchd-trade.err</string>
</dict>
</plist>
