# automation/com.trader.validate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent: monthly walk-forward re-validation (LOCAL).
  Runs on the 1st of each month at 09:00 local (Europe/London). yfinance works
  locally, unlike the cloud sandbox which blocks Yahoo Finance egress.
  Absolute paths required (launchd does not expand ~ / $HOME).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trader.validate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/snigdhatiwari/trader/automation/run_validate.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Day</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/snigdhatiwari/trader/logs/launchd-validate.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/snigdhatiwari/trader/logs/launchd-validate.err</string>
</dict>
</plist>
